# etc/NetworkManager/dispatcher.d/99-set-timezone
#!/usr/bin/env sh

#############################################################################################################################################################################
#   The license used for this file and its contents is: BSD-3-Clause                                                                                                        #
#                                                                                                                                                                           #
#   Copyright <2025> <Uri Herrera <uri_herrera@nxos.org>>                                                                                                                   #
#                                                                                                                                                                           #
#   Redistribution and use in source and binary forms, with or without modification, are permitted provided that the following conditions are met:                          #
#                                                                                                                                                                           #
#    1. Redistributions of source code must retain the above copyright notice, this list of conditions and the following disclaimer.                                        #
#                                                                                                                                                                           #
#    2. Redistributions in binary form must reproduce the above copyright notice, this list of conditions and the following disclaimer                                      #
#       in the documentation and/or other materials provided with the distribution.                                                                                         #
#                                                                                                                                                                           #
#    3. Neither the name of the copyright holder nor the names of its contributors may be used to endorse or promote products derived from this software                    #
#       without specific prior written permission.                                                                                                                          #
#                                                                                                                                                                           #
#    THIS SOFTWARE IS PROVIDED BY THE COPYRIGHT HOLDERS AND CONTRIBUTORS "AS IS" AND ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO,                      #
#    THE IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE ARE DISCLAIMED. IN NO EVENT SHALL THE COPYRIGHT HOLDER OR CONTRIBUTORS                  #
#    BE LIABLE FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE                 #
#    GOODS OR SERVICES; LOSS OF USE, DATA,   OR PROFITS; OR BUSINESS INTERRUPTION) HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT,                      #
#    STRICT LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF SUCH DAMAGE.   #
#############################################################################################################################################################################


# -- Verify commands exist.

for cmd in jq curl nm-online; do
    command -v "$cmd" >/dev/null 2>&1 || exit 1
done


# -- Dispatcher variables.
  # shellcheck disable=SC2034
IFACE="$1"
STATUS="$2"

[ "$STATUS" != "up" ] && exit 0


# -- Wait until full connectivity is established, including DNS.

nm-online -q -t 15 >/dev/null 2>&1 || exit 1


# -- Use IP geolocation to detect timezone.
#    -- ⚠️ A Note on Privacy: This method sends your IP address to a third-party service (worldtimeapi.org). This is a common practice and generally low-risk, but it's important to be aware of.

TZ_API="https://worldtimeapi.org/api/ip"
ZONE=$(curl -fsSL "$TZ_API" | jq -r .timezone)


# -- If successful, set it. Fallback using a cached timezone (from /etc/timezone).

if [ -n "$ZONE" ] && [ "$ZONE" != "null" ] && [ -e "/usr/share/zoneinfo/$ZONE" ]; then
    ln -sf "/usr/share/zoneinfo/$ZONE" /etc/localtime
    echo "$ZONE" > /etc/timezone
    echo "Timezone set to $ZONE"
else
    if [ -f /etc/timezone ] && [ -e "/usr/share/zoneinfo/$(cat /etc/timezone)" ]; then
        CACHED_ZONE="$(cat /etc/timezone)"
        ln -sf "/usr/share/zoneinfo/$CACHED_ZONE" /etc/localtime
        echo "Timezone fallback to $CACHED_ZONE"
    else
        echo "Failed to determine timezone automatically and no fallback available."
    fi
fi


# -- Restart NTP service to re-sync with the correct offset.

rc-service ntpsec-openrc restart
